# root/etc/services.d/exabgp/run
#!/bin/sh
exabgp --folder /etc/exabgp /etc/exabgp/exabgp.conf
